# skills/remind-me/systemd/remind-me.service
[Unit]
Description=Remind Me Skill Wakeup Handler
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.config/remind-me-skill/scripts/linux/wakeup_handler.sh
StandardOutput=journal
StandardError=journal